# 0x1A-application_server/gunicorn.service
[Unit]
Description=Gunicorn service for your_project_directory
After=network.target

[Service]
